# cmn/src/dist/ansible/roles/supervisor/templates/etc/supervisor/conf.d/program.conf.j2
{% for key, value in supervisor_programs.iteritems() %}
[program:{{key}}]
command={{value.command}}
user=ubuntu
stdout_logfile=/var/log/supervisor/{{key}}-stdout.log
stderr_logfile=/var/log/supervisor/{{key}}-stderr.log
environment={{value.env}}

{% endfor %}