# cmn/src/main/dist/ansible/roles/kibana/files/etc/systemd/system/kibana.service
[Unit]
Description=kibana
After=network.target

[Service]
User=kibana
ExecStart=/opt/kibana-4.6.1-linux-x86_64/bin/kibana

[Install]
WantedBy=multi-user.target
