# cmn/src/main/dist/ansible/roles/kafka/files/etc/systemd/system/zookeeper.service
[Unit]
Description=zookeeper
After=network.target

[Service]
User=kafka
ExecStart=/opt/kafka_2.12-2.2.0/bin/zookeeper-server-start.sh /opt/kafka_2.12-2.2.0/config/zookeeper.properties

[Install]
WantedBy=multi-user.target
