# cmn/src/main/dist/ansible/roles/kafka/templates/etc/systemd/system/kafka.service.j2
[Unit]
Description=kafka
After=network.target zookeeper.service

[Service]
User=kafka
Environment="KAFKA_HEAP_OPTS={{kafka_java_opts}}"
ExecStart=/opt/kafka_2.12-2.2.0/bin/kafka-server-start.sh /opt/kafka_2.12-2.2.0/config/server.properties

[Install]
WantedBy=multi-user.target
